# system/usr__lib__systemd__user__flatpak-user-update.service
[Unit]
Description=Update User Flatpaks

[Service]
Type=oneshot
ExecStart=/usr/bin/flatpak update --assumeyes --noninteractive

[Install]
WantedBy=default.target
