# dev-toolbox/base/etc/profile.d/podman-host.sh
# /etc/profile.d/podman-host.sh - sets Podman and Docker host.
export DOCKER_HOST="unix:///run/user/1000/podman/podman.sock"

